# supervisord.conf
[supervisord]
nodaemon=true

[program:squid]
command=/usr/sbin/squid -N -f /etc/squid/squid.conf
autorestart=true

[program:frpc]
command=/usr/local/bin/frpc -c /etc/frp/frpc.ini
autorestart=true